# deploy/systemd/opengram-web.service
[Unit]
Description=OpenGram Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opengram
Group=opengram
WorkingDirectory=/opt/opengram/web
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000
Environment=DATABASE_URL=/opt/opengram/data/opengram.db
Environment=OPENGRAM_CONFIG_PATH=/opt/opengram/config/opengram.config.json
EnvironmentFile=-/opt/opengram/config/opengram.env
ExecStart=/usr/bin/env node /opt/opengram/web/dist/server/server.mjs
Restart=always
RestartSec=2
MemoryHigh=256M
MemoryMax=384M
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/opengram/data /opt/opengram/config

[Install]
WantedBy=multi-user.target
